# pkg/greendot-helper.service
[Unit]
Description=GreenDotRDMA privileged helper
Documentation=https://github.com/asdacap/greendotrdma

[Service]
Type=simple
# Runs as root (uid 0) for configfs/PAM/zfs, but with group greendot so the
# RuntimeDirectory (/run/greendotrdma) is owned root:greendot — letting the
# unprivileged greendot-web traverse it to reach the helper socket.
Group=greendot
ExecStart=/usr/libexec/greendotrdma/greendot-helper
Restart=on-failure
RestartSec=2
# Must stay root for configfs/PAM/zfs, and must run in the HOST mount namespace:
# the helper mounts filesystems (`zfs mount`, the auto-mount of `zfs create`,
# btrfs resize mounts) that the host and greendot-web have to see. Any
# mount-namespace sandbox option — PrivateTmp=, ProtectHome=, ProtectSystem=,
# ReadOnlyPaths=/ReadWritePaths= — gives the unit its own namespace with `slave`
# propagation, so mounts it makes never reach the host: `zfs list` then shows the
# dataset unmounted while a retry `zfs mount` fails "filesystem already mounted".
# So do NOT add them. Likewise do NOT add ProtectKernelTunables=/
# ProtectControlGroups=/RestrictAddressFamilies=: the NFS apply path writes
# /proc/fs/nfsd/portlist and runs `systemctl start nfs-server`, which they would
# silently break.
RuntimeDirectory=greendotrdma
RuntimeDirectoryMode=0750
RuntimeDirectoryPreserve=yes

[Install]
WantedBy=multi-user.target
