# pkg/greendot-helper.service
[Unit]
Description=GreenDotRDMA privileged helper
Documentation=https://github.com/asdacap/greendotrdma

[Service]
Type=simple
# Runs as root (uid 0) for configfs/PAM/zfs, but with group greendot so the
# RuntimeDirectory (/run/greendotrdma) is owned root:greendot — letting the
# unprivileged greendot-web traverse it to reach the helper socket.
Group=greendot
ExecStart=/usr/libexec/greendotrdma/greendot-helper
Restart=on-failure
RestartSec=2
# Must stay root for configfs/PAM/zfs, but lock down what we can.
ProtectHome=yes
PrivateTmp=yes
RuntimeDirectory=greendotrdma
RuntimeDirectoryMode=0750
RuntimeDirectoryPreserve=yes

[Install]
WantedBy=multi-user.target
